# scripts/com.leesy.blog-image-audit.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>cc.leesy.blog-image-audit</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>-B</string>
    <string>-m</string>
    <string>scripts.audit_images</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/leesdove/Documents/blog</string>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/leesdove/Library/Logs/blog-image-audit.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/leesdove/Library/Logs/blog-image-audit.log</string>
</dict>
</plist>
